# mmega.service
[Unit]
Description=mmega uWSGI systemd service.

[Service]
User=
WorkingDirectory=
ExecStart=
Restart=on-failure

[Install]
WantedBy=multi-user.target
